# crontab/rpm-mirror
# Regular cron jobs for the rpm-mirror package
# Run every 6 hours

0 */6 * * *   rpm-mirror /usr/local/bin/rsync-repo.sh
